# systemd/travelrouter.service
[Unit]
Description=travelrouter management API and web UI
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/travelrouter
Environment=TROUTER_PORT=8080
ExecStart=/opt/travelrouter/venv/bin/python -m trouter.run
Restart=on-failure
RestartSec=3
# Hardening (relaxed where root networking control is required).
NoNewPrivileges=false
ProtectHome=false
ProtectSystem=full

[Install]
WantedBy=multi-user.target
